# autorun_scripts/Linux/dmonitor.service.template
[Unit]
Description=DMonitor
After=multi-user.target

[Service]
Type=simple
ExecStart=path

[Install]
WantedBy=multi-user.target